# com.quantum.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quantum.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/danielstevens/Desktop/github-projects/tools-automation/mcp_server.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/danielstevens/Library/Logs/tools-automation/mcp_server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/danielstevens/Library/Logs/tools-automation/mcp_server.err.log</string>
</dict>
</plist>
