# com.tools-automation.mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tools-automation.mcp-server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/danielstevens/Desktop/github-projects/tools-automation/mcp/servers/start_mcp_server.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/mcp-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/mcp-server-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/danielstevens/Desktop/github-projects/tools-automation/mcp/servers</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MCP_PORT</key>
        <string>5005</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
