# monitoring/templates/launchd-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.quantum.mcp</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/danielstevens/Desktop/Code/Tools/Automation/.venv/bin/python</string>
      <string>/Users/danielstevens/Desktop/Code/Tools/Automation/mcp_server.py</string>
      <string>127.0.0.1</string>
      <string>5005</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danielstevens/Desktop/Code/Tools</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/danielstevens/Desktop/Code/Tools/Automation/logs/mcp.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielstevens/Desktop/Code/Tools/Automation/logs/mcp.err</string>
  </dict>
</plist>
