# monitoring/templates/launchd-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.quantum.mcp-agent</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/danielstevens/Desktop/Code/Tools/Automation/.venv/bin/python</string>
      <string>/Users/danielstevens/Desktop/Code/Tools/Automation/agents/run_agent.py</string>
      <string>--name</string>
      <string>mcp-agent</string>
+      <string>--capabilities</string>
+      <string>execute</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danielstevens/Desktop/Code/Tools</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/danielstevens/Desktop/Code/Tools/Automation/logs/agent.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielstevens/Desktop/Code/Tools/Automation/logs/agent.err</string>
  </dict>
</plist>
