# monitoring/templates/systemd-agent.service
[Unit]
Description=MCP Agent
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/Users/danielstevens/Desktop/Code/Tools
ExecStart=/Users/danielstevens/Desktop/Code/Tools/Automation/.venv/bin/python Automation/agents/run_agent.py --name mcp-agent --capabilities execute
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
